# tools/DedicatedDaemonDebian/debian/postrm
#!/bin/sh -e
update-rc.d -f openlierox-dedicated remove || true
if [ "$1" = "purge" ] ; then 
	rm -rf /etc/openlierox-dedicated || true
	rm -rf /var/games/openlierox-dedicated || true
fi
exit 0
