# debian/postinst
#!/bin/sh -e
update-rc.d openlierox-dedicated defaults || true
exit 0
